# flaskapp.conf
start on runlevel [2345]
stop on runlevel [!2345]

setuid nitrous
setgid www-data

chdir /home/nitrous/flaskapp
env PATH=.env/bin
exec uwsgi .ini